# docker/supervisord.conf
[supervisord]
nodaemon=false

[program:discord-irc-billimek]
priority=30
directory=/discord-irc
command=discord-irc-billimek --config /discord-irc/config.json
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/discord-irc-billimek.log
stderr_logfile=/var/log/discord-irc-billimek.err

[unix_http_server]
file=%(here)s/supervisor.sock

[supervisorctl]
serverurl=unix://%(here)s/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
